# backend/setup/giftclick.conf
# Respawn up to 10 times in 1 minute.
respawn
respawn limit 10 1

# Start on all runlevels, except recovery (1). Stop on halt and reboot.
start on runlevel [2345]
stop on runlevel [06]

script
	# Run from geon's `~/giftclick`.
	cd /home/geon/giftclick

	# Run as geon. `exec` makes node replace the shell process. No point in leaving it running.
	# Append stdout and stderr to separate log files.
	exec sudo -u geon npm start 1>> ../debug-giftclick.log 2>> ../error-giftclick.log
end script
